# templates/default/mysqlsla.cron.erb
SHELL=/bin/sh
0 5 * * * root  /usr/local/bin/mysqlsla -lt general <%= @mysql['query_log'] %>.1 >> /var/log/mysqlsla-common.log
0 6 * * * root  /usr/local/bin/mysqlsla -lt slow <%= @mysql['slow_query_log'] %>.1 >> /var/log/mysqlsla-slow.log
